# config/includes.chroot_after_packages/usr/lib/systemd/system/oxyos-apt-setup.service
[Unit]
Description=oxyos apt setup service
After=getty@%i.service

[Service]
ExecStart=/usr/lib/oxyos/apt-setup
Type=oneshot

[Install]
WantedBy=multi-user.target
